# rootfs/etc/services.d/.s6-svscan/finish
#!@execline@/bin/execlineb -P
@execline@/bin/exit 0
