# rootfs/etc/services.d/cleantmp/run
#!@execline@/bin/execlineb -P
@execline@/bin/export PATH @execline@/bin:@s6PortableUtils@/bin
foreground { s6-sleep 10 }
foreground { @findutils@/bin/find /tmp/ -ctime +3 -delete }
s6-sleep 7200
